# ops/systemd/pacifica-full-fidelity-r2-lifecycle.timer
[Unit]
Description=Run Pacifica full-fidelity R2 lifecycle every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=pacifica-full-fidelity-r2-lifecycle.service

[Install]
WantedBy=timers.target
